# lib/opkg/info/glibc-opt.list
/opt/etc/init.d/rc.unslung
/opt/etc/init.d/rc.func
/opt/etc/profile
