# lib/opkg/info/minidlna.list
/opt/bin/minidlna
/opt/etc/minidlna.conf
/opt/etc/init.d/S90minidlna
